# home/dot_config/systemd/user/tmux-sync-theme.path
[Unit]
Description=Watch KDE color scheme changes and sync tmux theme

[Path]
PathChanged=%h/.config/kdeglobals
Unit=tmux-sync-theme.service

[Install]
WantedBy=default.target
